# systemd/pgzabbix-discover.timer
[Unit]
Description=pgZabbix discovery

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

