# systemd/pgzabbix.timer
[Unit]
Description=pgZabbix stat reader

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

